# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/home/ubuntu/peace_online/peace_online/peace_online.sock
[Install]
WantedBy=sockets.target